# scripts/duckgres.service
[Unit]
Description=Duckgres - PostgreSQL wire protocol server with DuckDB backend
After=network.target

[Service]
Type=notify
NotifyAccess=all
User=duckgres
Group=duckgres
WorkingDirectory=/opt/duckgres
# Ensure data directory exists and is owned by the service user
ExecStartPre=/bin/mkdir -p /data
ExecStartPre=/bin/chown duckgres:duckgres /data
ExecStart=/opt/duckgres/duckgres --config /opt/duckgres/duckgres.yaml --mode control-plane --socket-dir /var/run/duckgres --handover-socket /var/run/duckgres/handover.sock
# Graceful handover: SIGUSR1 triggers self-exec, the new process takes over
# listener FDs and workers via the handover socket, then does a rolling update.
ExecReload=/bin/kill -USR1 $MAINPID
Restart=always
RestartSec=10
LimitNOFILE=65535
# New CP must send READY=1 within this window after SIGUSR1 handover.
TimeoutStartSec=60
# Give control plane time to drain workers gracefully
TimeoutStopSec=60

# Set HOME for DuckDB extension storage (ProtectHome=true blocks /home)
Environment=HOME=/data

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/data

# systemd creates and manages /var/run/duckgres for Unix sockets.
# CRITICAL: RuntimeDirectoryPreserve=yes is required for the rolling-restart
# model. Without it, systemd may unmount or remount the directory as RO
# when the 'old' process exits during a handover, even after a new MAINPID
# has been registered.
RuntimeDirectory=duckgres
RuntimeDirectoryPreserve=yes
ReadWritePaths=/var/run/duckgres

[Install]
WantedBy=multi-user.target
